# domain_check/domain-check.ini
[program:domain-check]
command=python3 /usr/local/lib/python3.6/site-packages/domaincheck/DomainCheck.py
stdout_logfile=/var/log/domain-check_out.log
stderr_logfile=/var/log/domain-check_err.log
user=root
startsecs=8
startretries=5
autostart=true
autorestart=true
